# files/etc/cron.d/crontab
0 0 * * * /nightly.cron.sh >/dev/null 2>&1
